# wohnungssucher.service
[Unit]
Description=Wohnungssucher

[Service]
User=wsuser
Type=oneshot
ExecStart=/usr/bin/python3 /usr/bin/wohnungssucher/main.py

