# wohnungssucher.service
[Unit]
Description=Wohnungssucher

[Service]
User=wsuser
Type=oneshot
ExecStart=/usr/bin/python3 /opt/wohnungssucher/main.py

